# flask/panel-api.service
[Unit]
Description=LED Panel API

[Service]
ExecStart=/home/pi/pi/flask/panel-api

[Install]
WantedBy=multi-user.target
